# backend/etc/crontab
SHELL=/bin/sh
0 0 * * * /app/reset_tokens.sh >> /app/var/log/cron.log 2>&1
0 0 * * * /app/del_guests.sh >> /app/var/log/cron.log 2>&1
